# templates/launchd/com.wavex-os.resource-sweep.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wavex-os.resource-sweep</string>
  <key>ProgramArguments</key>
  <array>
    <string>${STATE_DIR}/resource-sweep.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>COMPANY_ID</key>
    <string>${COMPANY_ID}</string>
    <key>API_BASE</key>
    <string>${API_BASE}</string>
    <key>STATE_DIR</key>
    <string>${STATE_DIR}</string>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>${STATE_DIR}/resource-sweep.log</string>
  <key>StandardErrorPath</key>
  <string>${STATE_DIR}/resource-sweep.log</string>
</dict>
</plist>
